# docker-monitor.service
[Service]
WorkingDirectory=/opt/docker-monitor
ExecStart=/opt/docker-monitor/venv/bin/python monitor_and_alert.py
User=fu
Group=fu
StandardOutput=journal
StandardError=journal
Restart=always
